# crontab
0 * * * * root python /home/ubuntu/weather/fuzzy-meme/getweatherstats.py
30 * * * * root python /home/ubuntu/weather/fuzzy-meme/getweatherstats.py
